# dasting.service.template
[Unit]
Description=dastin

[Service]
ExecStart=/path/to/installed/bin
Restart=on-failure

EnvironmentFile=/path/to/env.conf
[Install]
WantedBy=multi-user.target